# dot_config/systemd/user/mako.service
[Unit]
Description=A lightweight Wayland notification daemon
Documentation=man:mako(1)
PartOf=sway-session.target
Requires=sway-session.target
After=sway-session.target

[Service]
Type=dbus
ExecStart=/usr/bin/mako
ExecReload=/usr/bin/makoctl reload
BusName=org.freedesktop.Notifications
Slice=session.slice

[Install]
WantedBy=sway-session.target
